# deployment/systemd/mollo-vault-import.timer
[Unit]
Description=Mollo Vault Import — every 15 min

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true
Unit=mollo-vault-import.service

[Install]
WantedBy=timers.target
